# srcs/db/tools/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:mysql]
command=mysqld -u root --console
priority=2

[program:db]
command=/db.sh
